# 17g.initd
#!/sbin/openrc-run
 # Copyright 1999-2013 Gentoo Foundation
 # Distributed under the terms of the GNU General Public License v2
 # $Header: $

start() {
    /lib/live-install/service.py
}
